# bin/sandlights.service
[Unit]
Description=SandTable Lights Service

[Service]
WorkingDirectory=/home/dubnom/sandtable
ExecStart=/home/dubnom/.venv/bin/python src/ledaemon.py
Restart=always
RestartSec=10
SyslogIdentifier=sandtable-lights
AmbientCapabilities=CAP_NET_BIND_SERVICE
User=sandtable

[Install]
WantedBy=multi-user.target

